# mongodb/init-replicaset.conf
[program:init-replicaset]
command = sh -c 'sleep 10 && /usr/bin/mongo localhost:27017 </tmp/init.js'
autorestart = false
