# booktype/booktype.okfn.org.conf
[program:booktype]
command = gunicorn_django -b 0.0.0.0:8000 -w 4
directory = /
autorestart = true
autostart = true
